# rootfs/etc/services.d/crowdsec-dashboard/finish
#!/usr/bin/with-contenv bashio
bashio::log.info "CrowdSec Dashboard stopped."
