# rootfs/etc/services.d/crowdsec-dashboard/run
#!/usr/bin/with-contenv bashio

# Read options from HA addon config
CROWDSEC_URL=$(bashio::config 'crowdsec_url')
MACHINE_ID=$(bashio::config 'machine_id')
MACHINE_PASSWORD=$(bashio::config 'machine_password')
LOG_LEVEL=$(bashio::config 'log_level')

export CROWDSEC_URL
export MACHINE_ID
export MACHINE_PASSWORD
export LOG_LEVEL

bashio::log.info "Starting CrowdSec Dashboard..."
bashio::log.info "CrowdSec LAPI: ${CROWDSEC_URL}"
bashio::log.info "Machine ID: ${MACHINE_ID}"

exec python3 /app/app.py
